# config/piheating.service
[Unit]
Description=Pi Heating
After=network.target
[Service]
ExecStart=/bin/bash /home/pi/PiHeating/startupPiHeating.sh
WorkingDirectory=/run/shm/PiHeating//bin
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
